# .config/systemd/user/yadm-auto-updater.timer
[Unit]
Description=Update yadm dotfiles periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
