# .config/systemd/user/yadm-auto-updater.service
[Unit]
Description=Auto pull yadm

[Service]
Type=oneshot
ExecStart=/usr/bin/yadm pull
WorkingDirectory=%h
